# supervisor/healthcheckd.ini
[program:healthcheckd]

command=/usr/bin/python /usr/local/bin/healthcheckd.py
autostart=true
autorestart=true
