# supervisor/puppet-agent.ini
[program:puppet-agent]

command=/bin/bash /usr/local/src/puppet-masterless/localpuppetmaster.sh -D -d /usr/local/src/localpuppetmaster/master -s /usr/local/src/localpuppetmaster/saltmaster.pp
autostart=true
autorestart=true
