# supervisor/saltapi.ini
[program:saltapi]

command=/usr/bin/salt-api
autostart=true
autorestart=true
